# infra/dispatch-ui.service
[Unit]
Description=Dispatch UI (vite preview, static build)
Documentation=https://github.com/MrSchrodingers/debt-adb-framework
After=network-online.target dispatch-core.service
Wants=network-online.target

[Service]
Type=simple
User=adb
WorkingDirectory=/var/www/debt-adb-framework/packages/ui
Environment=NODE_ENV=production
Environment=NVM_DIR=/home/adb/.nvm
Environment=PATH=/home/adb/.nvm/versions/node/v22.22.2/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

ExecStartPre=/bin/bash -c 'test -f /var/www/debt-adb-framework/packages/ui/dist/index.html'
ExecStart=/home/adb/.nvm/versions/node/v22.22.2/bin/node /var/www/debt-adb-framework/packages/ui/node_modules/vite/bin/vite.js preview --host 127.0.0.1 --port 5174

Restart=always
RestartSec=5
StartLimitInterval=300
StartLimitBurst=10

NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/www/debt-adb-framework/packages/ui
ProtectHome=read-only
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

LimitNOFILE=65536
MemoryHigh=512M
MemoryMax=1G

StandardOutput=journal
StandardError=journal
SyslogIdentifier=dispatch-ui

[Install]
WantedBy=multi-user.target
